# embedded-linux-stm32mp1-labs/thirdparty/target/etc/inittab
::sysinit:/etc/init.d/rcS
ttySTM0::askfirst:/bin/sh
::restart:/sbin/init
